# peer/config/supervisor/peer-server.conf
[program:freetime-peer-1]
; FreeTime Peer Server Process Control Configuration
; Supervisor will monitor and auto-restart peer processes

; Program identifier
command=/usr/bin/node /home/freetime/peer/peer-server.js
directory=/home/freetime/peer
user=freetime
group=freetime

; Logging
stdout_logfile=/var/log/freetime-peer.log
stdout_logfile_maxbytes=10485760
stdout_logfile_backups=10
stderr_logfile=/var/log/freetime-peer-error.log
stderr_logfile_maxbytes=10485760
stderr_logfile_backups=10

; Process management
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true

; Resource limits
numprocs=1
priority=999
umask=0022

; Environment
environment=NODE_ENV="production",PATH="/usr/local/bin:/usr/bin:/bin"

; Attempt to catch OOM signals
; supervisor will restart process if memory grows too large
; Adjust memmon_available_mb based on system RAM
; memmon_available_mb=4096

; Health checks
; Check process is actually listening
; If not, restart it
eventlistener=memmon
buffer_size=100
events=TICK_60
result_handler=unix_listener
action_type=escalate
priority=998

; Graceful shutdown
stopwaitsecs=30

; Process restart policy
; fatal: Do not restart if exits with code != 0
; Process will be marked as fatal and logged
